# .config/autostart/__synccolor.desktop
[Desktop Entry]
Exec=__synccolor
Icon=application-x-shellscript
Name=__synccolor
Type=Application
X-KDE-AutostartScript=true
