# .config/autostart/checkbg-watch-awake.desktop
[Desktop Entry]
Exec=__checkbg_watch_awake
Icon=application-x-shellscript
Name=__checkbg_watch_awake
Type=Application
X-KDE-AutostartScript=true
