# .config/autostart/checkbg-periodic.desktop
[Desktop Entry]
Exec=__checkbg_periodic
Icon=application-x-shellscript
Name=__checkbg_periodic
Type=Application
X-KDE-AutostartScript=true
