# .config/autostart/checkbg.desktop
[Desktop Entry]
Exec=checkbg
Icon=application-x-shellscript
Name=checkbg
Type=Application
X-KDE-AutostartScript=true
