# .config/autostart/watch-awake-checkbg.desktop
[Desktop Entry]
Exec=watch-awake-checkbg
Icon=application-x-shellscript
Name=watch-awake-checkbg
Type=Application
X-KDE-AutostartScript=true
